# aio/s6-rc.d/postgres-bootstrap/up
#!/command/with-contenv bash
exec /scripts/bash/init-db.sh bootstrap
