# aio/s6-rc.d/postgres-init/up
#!/bin/sh
exec /scripts/bash/init-db.sh provision
