# aio/s6-rc.d/process-metrics/finish
#!/command/with-contenv bash
set -eu

STATE_DIR="/run/process-metrics"
COUNT_FILE="$STATE_DIR/restart-count"
MAX_RESTARTS="${MAX_RESTARTS:-3}"

mkdir -p "$STATE_DIR"

count=$(cat "$COUNT_FILE" 2>/dev/null || echo 0)
count=$((count+1))
echo "$count" > "$COUNT_FILE"
date +%s > "$STATE_DIR/last-restart"

echo "[process-metrics] crash count: $count/${MAX_RESTARTS}"

if [ "$count" -ge "$MAX_RESTARTS" ]; then
  echo "[process-metrics] too many failures, terminating supervision tree"
  exec /command/s6-svscanctl -t /run/service
fi

exit 0
