# aio/s6-rc.d/postgres/finish
#!/command/with-contenv bash
set -eu

STATE_DIR="/run/postgres"
COUNT_FILE="$STATE_DIR/restart-count"
MAX_RESTARTS="${MAX_RESTARTS:-3}"

mkdir -p "$STATE_DIR"

# Count every restart (normal or abnormal)
count=$(cat "$COUNT_FILE" 2>/dev/null || echo 0)
count=$((count+1))
echo -n "$count" > "$COUNT_FILE"

echo "[postgres] crash count: $count/${MAX_RESTARTS}"

if [ "$count" -ge "$MAX_RESTARTS" ]; then
  echo "[postgres] too many failures, terminating supervision tree"
  # Ask s6-svscan to exit; this will stop the container
  exec /command/s6-svscanctl -t /run/service
fi

exit 0
