# aio/s6-rc.d/postgres/run
#!/bin/sh
exec su -s /bin/bash -c "postgres -D ${PGDATA:-/sp/postgresql/data}" postgres 2>&1 | tee /var/log/postgresql.log
