# aio/s6-rc.d/postgres/run
#!/bin/sh
exec su -s /bin/bash -c "postgres -D ${PGDATA:-/var/lib/postgresql/data}" postgres
